# crontab
0 12 * * fri,sun /bin/auto_reminder --api --config /etc/auto_reminder/config.json
